# root/etc/services.d/php/run
#!/usr/bin/with-contenv bash

export HOME="~app"
exec s6-setuidgid root "/conf/start_php.sh"
